# freebsd/crontab
MAILTO=root

0 * * * *       root    /var/scripts/snapshot-hourly.sh

0 0 * * *       root    /var/scripts/snapshot-daily.sh

*/15 * * * *    root    /var/scripts/check.sh
